# templates/default/fakes3.upstart.conf.erb
#!upstart
description "fakes3 server"
author      "pluma"

start on startup
stop on shutdown

# Automatically restart process if crashed
respawn

# Run before process
pre-start script
end script

# Start the process
script
   /usr/local/bin/fakes3 -r /srv/s3 -p 80 --limit=50K
end script